# ansible/files/tmpclean.crontab
* 3 * * * root /bin/tmpwatch 336 /opt/dlp-selfdeposit/current/tmp/imports
